# ambientLightFix_intel.service
[Unit]
Description=Automatically adjust screen-brightness

[Service]
Type=oneshot
ExecStart=/usr/bin/ambientLightFix_intel

[Install]
WantedBy=multi-user.target
