# ambientLightFix_gmux.service
[Unit]
Description=Automatically adjust screen-brightness

[Service]
Type=oneshot
ExecStart=/usr/bin/ambientLightFix_gmux

[Install]
WantedBy=multi-user.target
