# ambientLightFix_radeon.service
[Unit]
Description=Automatically adjust screen-brightness

[Service]
Type=oneshot
ExecStart=/usr/bin/ambientLightFix_radeon

[Install]
WantedBy=multi-user.target
